# examples/htcondor/files/condor.ini
[program:condor]
command=/usr/local/sbin/condor-wrapper.sh
